# deploy/supervisord_native.conf
; ================================================================================
;  CONFIGURACIN SUPERVISOR PARA PRODUCCIN NATIVA
; ================================================================================
;
; Ubicaci贸n: /etc/supervisor/conf.d/telegram360.conf
; 
; Crear archivo:
; $ sudo nano /etc/supervisor/conf.d/telegram360.conf
; 
; Copiar TODO el contenido de este archivo
; 
; Luego:
; $ sudo supervisorctl reread
; $ sudo supervisorctl update
; $ sudo supervisorctl status
;
; ================================================================================

; 
; [program:telegram-backend]
; 
; Servicio: Backend FastAPI con Gunicorn
; Puerto: 8000
; Descripci贸n: API REST para la aplicaci贸n
; 

[program:telegram-backend]

; Comando a ejecutar
; Usando Gunicorn con 2 workers (uno por CPU)
command=/home/telegramapp/app/venv/bin/gunicorn \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --access-logfile /home/telegramapp/app/logs/access.log \
    --error-logfile /home/telegramapp/app/logs/error.log \
    -m 007 \
    plataforma_canales_back-main.routes:app

; Directorio de trabajo
directory=/home/telegramapp/app

; Usuario a ejecutar como
user=telegramapp

; Archivo PID
pidfile=/home/telegramapp/app/logs/backend.pid

; 
; AUTO-START y AUTO-RESTART
; 

; Iniciar autom谩ticamente cuando arranca el servidor
autostart=true

; Reiniciar autom谩ticamente si se detiene
autorestart=true

; Esperar a que se inicie completamente
startsecs=10

; Intentar reiniciar 3 veces
startretries=3

; 
; LOGGING
; 

; Redirigir stderr a stdout
redirect_stderr=true

; Archivo de log de stdout
stdout_logfile=/home/telegramapp/app/logs/backend.log

; Tama帽o m谩ximo de log antes de rotar
stdout_logfile_maxbytes=50MB

; N煤mero de rotaciones a mantener
stdout_logfile_backups=5

; 
; SEALES Y PARADA
; 

; Se帽al para parar gracefully
stopsignal=TERM

; Tiempo de espera para parada graceful antes de KILL
stopwaitsecs=30

; Permitir que se escriba en stderr
silent=false

; 
; VARIABLES DE ENTORNO
; 

; Cargar .env.prod
environment=PATH="/home/telegramapp/app/venv/bin",PYTHONUNBUFFERED="1"

; 
; [program:telegram-workers]
;
; Servicio: Workers de Telegram (Telethon + Metrics)
; Descripci贸n: Procesa m茅tricas y eventos de canales Telegram
; 

[program:telegram-workers]

; Comando a ejecutar - Script de workers
command=/home/telegramapp/app/venv/bin/python \
    -u \
    telegram-workers/workers/metrics_tracker/metrics_tracker.py

; Directorio de trabajo
directory=/home/telegramapp/app

; Usuario a ejecutar como
user=telegramapp

; Archivo PID
pidfile=/home/telegramapp/app/logs/workers.pid

; 
; AUTO-START y AUTO-RESTART
; 

; Iniciar autom谩ticamente
autostart=true

; Reiniciar si se detiene
autorestart=true

; Esperar a que se inicie
startsecs=5

; Reintentos
startretries=3

; 
; LOGGING
; 

; Redirigir stderr a stdout
redirect_stderr=true

; Archivo de log
stdout_logfile=/home/telegramapp/app/logs/workers.log

; Tama帽o m谩ximo de log
stdout_logfile_maxbytes=50MB

; Rotaciones a mantener
stdout_logfile_backups=5

; 
; PARADA
; 

; Se帽al de parada
stopsignal=TERM

; Espera antes de KILL
stopwaitsecs=30

; 
; VARIABLES DE ENTORNO
; 

; Path y no buffering
environment=PATH="/home/telegramapp/app/venv/bin",PYTHONUNBUFFERED="1"

; 
; [group:telegram]
;
; Grupo de procesos para control conjunto
; 

[group:telegram]

; Programas incluidos en el grupo
programs=telegram-backend,telegram-workers

; Prioridad (orden de inicio)
priority=999

; 
; COMANDOS TILES PARA SUPERVISOR
; 

; Ver estado de todos
; sudo supervisorctl status

; Ver estado del grupo
; sudo supervisorctl status telegram

; Ver estado espec铆fico
; sudo supervisorctl status telegram:telegram-backend
; sudo supervisorctl status telegram:telegram-workers

; Iniciar todo
; sudo supervisorctl start telegram

; Parar todo
; sudo supervisorctl stop telegram

; Reiniciar todo
; sudo supervisorctl restart telegram

; Reiniciar espec铆fico
; sudo supervisorctl restart telegram:telegram-backend

; Ver logs en tiempo real
; sudo tail -f /home/telegramapp/app/logs/backend.log
; sudo tail -f /home/telegramapp/app/logs/workers.log

; Recargar configuraci贸n
; sudo supervisorctl reread
; sudo supervisorctl update

; Ejecutar comando en todos
; sudo supervisorctl start all
; sudo supervisorctl stop all
; sudo supervisorctl restart all

; 
; MONITOREO Y TROUBLESHOOTING
; 

; Ver procesos corriendo
; ps aux | grep -E "gunicorn|python" | grep -v grep

; Ver puertos escuchando
; sudo netstat -plnt | grep 8000

; Ver memoria usada
; ps aux --sort=-%mem | head -10

; Ver uso de CPU
; ps aux --sort=-%cpu | head -10

; Ver logs de error de Supervisor
; sudo tail -f /var/log/supervisor/supervisord.log

; Reiniciar supervisord completo
; sudo systemctl restart supervisor

; 
; CONFIGURACIN AVANZADA (Opcional)
; 

; Descomenta si quieres logging m谩s detallado de supervisor
; [supervisord]
; logfile=/var/log/supervisor/supervisord.log
; loglevel=debug

; Para RabbitMQ (si lo necesitas despu茅s)
; [program:rabbitmq]
; command=/usr/sbin/rabbitmq-server
; autostart=true
; autorestart=true
; stdout_logfile=/var/log/supervisor/rabbitmq.log
; redirect_stderr=true

; 
; PASOS PARA ACTIVAR
; 

; 1. Crear archivo:
;    $ sudo nano /etc/supervisor/conf.d/telegram360.conf

; 2. Copiar TODO el contenido de este archivo

; 3. Guardar (Ctrl+X, Y, Enter)

; 4. Recargar supervisor:
;    $ sudo supervisorctl reread
;    $ sudo supervisorctl update

; 5. Verificar estado:
;    $ sudo supervisorctl status
;    Esperado: telegram:telegram-backend RUNNING y telegram:telegram-workers RUNNING

; 6. Ver logs para debugging:
;    $ sudo tail -100f /home/telegramapp/app/logs/backend.log

; 
